# docker/crontab
PATH=/www/wwwroot/proxypanel
SHELL=/bin/bash

* * * * * sudo -u "www-data" /usr/local/bin/php artisan schedule:run > /proc/1/fd/1 2>&1
